# lib/epmd/priv/erlpmd.socket
[Unit]
Description=ErlPmd Server Activation Socket

[Socket]
ListenStream=127.0.0.1:4369
Accept=false

[Install]
WantedBy=sockets.target
